# stock-watchdog.path
[Unit]
Description=Watch for changes in stock-watchdog config file
Wants=stock-watchdog.service

[Path]
PathChanged=/opt/stock-watchdog/config.json
Unit=stock-watchdog.service

[Install]
WantedBy=multi-user.target 